# systemd/azazel-edge-epd-refresh.service
[Unit]
Description=Azazel-Edge EPD mode refresh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/azazel-edge-epd-refresh
TimeoutStartSec=60
StandardOutput=journal
StandardError=journal
