# systemd/azazel-edge-epd-refresh.timer
[Unit]
Description=Run Azazel-Edge EPD refresh periodically

[Timer]
OnBootSec=30s
OnUnitActiveSec=2min
Unit=azazel-edge-epd-refresh.service

[Install]
WantedBy=timers.target
